# wordpress-0.2.0/Templates/default/rc.local.erb
#!/bin/bash
#
# Generated by Chef - Do not modify as changes will be overridden next chef run.
# To add more lines, adjust the node attribute node[:melbit_linux][:rc_local_commands]

touch /var/lock/subsys/local

if [ ! -d /root/.ssh ] ; then
	mkdir -p /root/.ssh
        chmod 0700 /root/.ssh
        restorecon /root/.ssh
fi

# rhel bug zilla #707364
if [ ! -f /etc/blkid/blkid.tab ] ; then
	blkid /dev/xvda &>/dev/null
fi

<% @rc_local_commands.each do |rccommand| -%>
<%= rccommand %>
<% end -%>

